# backend/crontab
# NAV Sync — twice daily (01:00 & 13:00 UTC = 6:30 AM & 6:30 PM IST)
0 1 * * * /usr/local/bin/python /app/scripts/sync_amfi.py >> /proc/1/fd/1 2>&1
0 13 * * * /usr/local/bin/python /app/scripts/sync_amfi.py >> /proc/1/fd/1 2>&1

# NAV Gap Healer — once daily (03:30 UTC = 9:00 AM IST)
30 3 * * * /usr/local/bin/python /app/scripts/heal_nav_gaps.py >> /proc/1/fd/1 2>&1
